# templates/jupyterhub.service.epp
[Unit]
Description=Jupyterhub
After=network-online.target

[Service]
User=jupyterhub
Group=jupyterhub
Environment=PATH=/bin:/usr/bin:<%= $prefix %>/bin:<%= $slurm_home %>/bin
ExecStart=<%= $prefix %>/bin/jupyterhub --config /etc/jupyterhub/jupyterhub_config.json
WorkingDirectory=/run/jupyterhub

[Install]
WantedBy=multi-user.target
